# supervisord.conf
[supervisord]
nodaemon=true

[program:app-uwsgi]
command = uwsgi --ini /var/www/django/ini/uwsgi.ini

[program:nginx-app]
command = nginx